# packaging/deb/postinst
#!/bin/sh
set -e

systemctl daemon-reload || true
systemctl enable moblink || true
systemctl start moblink || true

exit 0
